# crontab

*/2 * * * *  python3 /sondensucher/Python/getradiosondycsv.py  >/dev/null 2>&1
